# rootfs/etc/cron.d/mixpost_crontab
* * * * * ( su -s /bin/sh -c 'cd /var/www/html && /usr/bin/php artisan schedule:run' www-data ) > /proc/1/fd/2 2>&1
